# logrotate/templates/logrotate.timer.j2
[Unit]
Description=Run daily logrotate

[Timer]
OnBootSec=15min
OnCalendar=daily
Persistent=true


[Install]
WantedBy=timers.target
